# deploy/once/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/sh -lc "/opt/venv/bin/trading-bot serve-api --host ${TRADING_BOT_API_HOST:-127.0.0.1} --port ${TRADING_BOT_API_PORT:-8000}"
directory=/app/backend
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:frontend]
command=/bin/sh -lc "./node_modules/.bin/srvx serve --prod --entry dist/server/server.js --host ${HOST:-0.0.0.0} --port ${PORT:-80}"
directory=/app/frontend
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
