# deploy/once/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
command=/opt/venv/bin/trading-bot serve-api --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:frontend]
command=/bin/sh -lc "./node_modules/.bin/srvx serve --prod --entry dist/server/server.js --host 0.0.0.0 --port 80"
directory=/app/frontend
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
